# NEMS/tests/nightly/cron/wcoss-user.cron
# NOTE: Make sure you change this to tide/luna/mars when that side is dev.

03 00 * * * mkdir -p /stmpd1/Ratko.Vasic /stmpd2/Ratko.Vasic /stmpd3/Ratko.Vasic
03 00 * * * mkdir -p /stmpp1/Ratko.Vasic /stmpp2/Ratko.Vasic
03 00 * * * mkdir -p /ptmpd1/Ratko.Vasic /ptmpd2/Ratko.Vasic /ptmpd3/Ratko.Vasic
03 00 * * * mkdir -p /ptmpp1/Ratko.Vasic /ptmpp2/Ratko.Vasic

03 00 * * * mkdir -p /gpfs/hps3/ptmp/Ratko.Vasic

33 16 * * *           /gpfs/hps3/emc/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-wcoss.sh ngt wcoss1         checkout > /ptmpp1/Ratko.Vasic/checkout-wcoss-phase1.log  2>&1
33 16 * * * ssh g20a1 /gpfs/hps3/emc/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-wcoss.sh ngt wcoss2         checkout > /ptmpp1/Ratko.Vasic/checkout-wcoss-phase2.log  2>&1
33 16 * * * ssh surge /gpfs/hps3/emc/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-wcoss.sh ngt wcoss_cray     checkout > /gpfs/hps3/ptmp/Ratko.Vasic/checkout-wcoss-cray.log  2>&1
33 16 * * * ssh venus /gpfs/hps3/emc/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-wcoss.sh ngt wcoss_dell_p3  checkout > /gpfs/hps3/ptmp/Ratko.Vasic/checkout-wcoss-dell-p3.log  2>&1

57 * * * * /u/Ratko.Vasic/RT/transfer.sh > /u/Ratko.Vasic/RT/transfer.log 2>&1
